# ops/supervisor/magento-consumers.conf
[program:erp-order-sync]
command=/usr/bin/php8.4 /home/user/htdocs/srv1113343.hstgr.cloud/bin/magento queue:consumers:start erp.order.sync.consumer --max-messages=1000
directory=/home/user/htdocs/srv1113343.hstgr.cloud
user=www-data
autostart=false
autorestart=false
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/home/user/htdocs/srv1113343.hstgr.cloud/var/log/erp-order-sync.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/home/user/htdocs/srv1113343.hstgr.cloud/var/log/erp-order-sync-error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3

[program:erp-order-sync-retry]
command=/usr/bin/php8.4 /home/user/htdocs/srv1113343.hstgr.cloud/bin/magento queue:consumers:start erp.order.sync.retry.consumer --max-messages=500
directory=/home/user/htdocs/srv1113343.hstgr.cloud
user=www-data
autostart=false
autorestart=false
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/home/user/htdocs/srv1113343.hstgr.cloud/var/log/erp-order-sync-retry.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/home/user/htdocs/srv1113343.hstgr.cloud/var/log/erp-order-sync-retry-error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3

[group:magento-erp]
programs=erp-order-sync,erp-order-sync-retry
priority=999
